# systemd/llama-proxy.service
[Unit]
Description=llama-proxy (role rewrite + thinking control, port 8000→8001)
After=network.target llama-server.service
Requires=llama-server.service

[Service]
Type=simple
User=YOUR_USERNAME
ExecStart=/path/to/python3 /opt/llama.cpp/llama-proxy.py
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/log/llama-proxy.log
StandardError=append:/var/log/llama-proxy.log

[Install]
WantedBy=multi-user.target
