# systemd/llama-server.service
[Unit]
Description=llama.cpp server (Qwen3-Coder-Next, Strix Halo)
After=network.target
Before=llama-proxy.service

[Service]
Type=simple
User=YOUR_USERNAME
ExecStart=/usr/local/bin/llama-server \
    --model /srv/ai/models/qwen3-coder-next/Qwen3-Coder-Next-Q4_K_M \
    --ctx-size 131072 \
    --parallel 1 \
    --host 127.0.0.1 \
    --port 8001 \
    -ngl 999 \
    -fa on \
    -dio \
    --jinja
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/llama-server.log
StandardError=append:/var/log/llama-server.log
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
